# samproxy.upstart
# Upstart job for samproxy, the Honeycomb Trace-Aware Sampling Proxy
# https://honeycomb.io/

description     "Samproxy Daemon"
author          "Ben Hartshorne <ben@honeycomb.io>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec su -s /bin/sh -c 'exec "$0" "$@"' honeycomb -- /usr/bin/samproxy -c /etc/samproxy/samproxy.toml -r /etc/samproxy/rules.toml
